# recipes-myapp/myproject-app/files/myproject-app.service
[Unit]
Description=MyProject application service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/myproject-app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
